# com.ocrmac.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ocrmac.api</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhanglixing/Dev/projects/250705_OcrMacApi/startup.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/zhanglixing/Dev/projects/250705_OcrMacApi</string>
    
    <key>StandardOutPath</key>
    <string>/Users/zhanglixing/Dev/projects/250705_OcrMacApi/logs/ocrmac_api.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/zhanglixing/Dev/projects/250705_OcrMacApi/logs/ocrmac_api_error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/zhanglixing/Dev/projects/250705_OcrMacApi:/Users/zhanglixing/Dev/projects/250705_OcrMacApi/src:/Users/zhanglixing/Dev/projects/250705_OcrMacApi/ocrmac-main</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist> 